# Base/op-node.service
[Unit]
Description=Optimism Node (Consensus)
After=network-online.target

[Service]
User=optimism
Group=optimism
Environment=HOME=/data/base/.op-node
EnvironmentFile=/etc/environment
Type=simple
ExecStart=/usr/local/bin/op-node
LimitNOFILE=8192
KillMode=process
KillSignal=SIGINT
TimeoutStopSec=90
Restart=on-failure
RestartSec=10s
#StandardOutput=null
#StandardError=null
SyslogIdentifier=op-node

[Install]
WantedBy=multi-user.target
